# deploy/systemd/ckp.service
[Unit]
Description=Cloud Knowledge Platform backend
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=ckp
Group=ckp
WorkingDirectory=/opt/ckp
EnvironmentFile=/opt/ckp/.env
ExecStart=/opt/ckp/.venv/bin/python -m uvicorn app.main:app --app-dir backend --host 127.0.0.1 --port 8787
Restart=on-failure
RestartSec=3
# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/ckp/vaults /opt/ckp/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
